# config/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;error_log /dev/stderr info;"
autostart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:odoo]
command=./launchodoo.sh
directory=/odoo/godoo_workspace
environment=PYTHONUNBUFFERED=1
autostart=true
autorestart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stder
redirect_stderr=true
stderr_logfile_maxbytes=0
